# ch12/crontab
*/1 * * * * root /usr/bin/sync.sh >> /tmp/log/cron.log 2>&1 
